# deploy/services.d/prometheus/run
#!/bin/env bash

prometheus --config.file=/etc/prometheus/prometheus.yml --storage.tsdb.path /data/prometheus
